# misc/gpsd.service
[Unit]
Description=gpsd startup

[Service]
Type=forking
ExecStart=gpsd /dev/ttyS2 -F /var/run/gpsd.sock

[Install]
WantedBy=multi-user.target

